# base/usr/lib/systemd/system/lightdm.service.d/steamos-autorecover.conf
[Unit]
OnFailure=steamos-autorepair.service
