# base/DEBIAN/steamos-lts-base.steamos-autorepair.service
[Unit]
Description=SteamOS Autorepair

[Service]
ExecStart=/usr/bin/steamos-autorepair.sh
Type=oneshot



